# docker/supervisord.conf
###############################################################################
# AFT supervisord program definitions.
#
# This file is included by the system /etc/supervisor/supervisord.conf via the
# `[include]` block, so we only define [program:*] sections here.
###############################################################################

[program:postgres]
command=/usr/lib/postgresql/18/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=30
priority=10
stdout_logfile=/var/log/aft/postgres.log
stderr_logfile=/var/log/aft/postgres.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:bun]
command=/usr/local/bin/bun run /app/index.ts
directory=/app
user=root
autostart=true
autorestart=true
startsecs=2
startretries=10
stopwaitsecs=15
priority=20
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    AFT_PROXY_SHARED_SECRET="%(ENV_AFT_PROXY_SHARED_SECRET)s",
    AFT_ADMIN_BOOTSTRAP_PASSWORD="%(ENV_AFT_ADMIN_BOOTSTRAP_PASSWORD)s",
    AFT_ADMIN_BOOTSTRAP_EMAIL="%(ENV_AFT_ADMIN_BOOTSTRAP_EMAIL)s",
    SMTP_HOST="%(ENV_SMTP_HOST)s",
    SMTP_PORT="%(ENV_SMTP_PORT)s",
    SMTP_SECURE="%(ENV_SMTP_SECURE)s",
    SMTP_USER="%(ENV_SMTP_USER)s",
    SMTP_PASS="%(ENV_SMTP_PASS)s",
    SMTP_FROM="%(ENV_SMTP_FROM)s",
    SMTP_REPLY_TO="%(ENV_SMTP_REPLY_TO)s"
stdout_logfile=/var/log/aft/bun.log
stderr_logfile=/var/log/aft/bun.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
stopsignal=QUIT
priority=30
stdout_logfile=/var/log/aft/nginx.log
stderr_logfile=/var/log/aft/nginx.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
